# 90-disk_check.rules
SUBSYSTEM=="block", RUN+="/usr/bin/bash -c '/usr/bin/kill -SIGUSR1 $(/usr/bin/cat /tmp/disk_check_pid)'"